# Services/ambient-ble-gateway.service
[Unit]
Description=Ambient Node BLE Gateway Service
After=network.target bluetooth.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ambient-node/Services
ExecStart=/home/pi/ambient-node/Services/.venv/bin/python3.11 /home/pi/ambient-node/Services/ble_gateway.py
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ambient-ble-gateway

[Install]
WantedBy=multi-user.target
