# BLE_gateway/ambient-ble-gateway.service
[Unit]
Description=Ambient Node BLE Gateway Service
After=network.target bluetooth.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ambient-node/rpi
ExecStart=/usr/bin/python3 /home/pi/ambient-node/rpi/ble_gateway.py
Restart=always
RestartSec=10

# Environment
Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ambient-ble-gateway

[Install]
WantedBy=multi-user.target
